# systemd/build-service.service
[Unit]
Description=Build Service
After=local-fs.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/build-service --config /etc/build-service/config.toml
Restart=on-failure
User=root

# Optional hardening (tune as needed):
# NoNewPrivileges=true
# PrivateTmp=true

# Optional runtime directory (if you set socket_path under /run/build-service/):
# RuntimeDirectory=build-service
# RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
